# files/etc/supervisor/conf.d/nginx-reloader.conf
[program:nginxreloader]
command=/usr/local/bin/nginxreloader
stdout_events_enabled=true
stderr_events_enabled=true
